# email/.config/systemd/user/mbsync.service
[Unit]
Description=Mailbox Synchronization Service

[Service]
Type=oneshot
ExecStart=/usr/bin/mbsync -V -c %h/.config/isync/mbsyncrc university
